# packages/etc/crontab
# Fetch and prepare data for mail at packages.debian.org
45 17 * * *	/org/packages.debian.org/bin/daily
@reboot         rm /org/packages.debian.org/files/cron.lock
