# packages/etc/crontab
# Fetch and prepare data for mail at packages.debian.org
55 2 * * *	/org/packages.debian.org/bin/daily
@reboot         rm -f /org/packages.debian.org/files/cron.lock
