# packages/etc/crontab
# Fetch and prepare data for mail at packages.debian.org
45 14 * * *	/org/packages.debian.org/bin/daily
